# mcsrv-timer-2.service
[Unit]
Description=McSrv Script Timer 2 Service

[Service]
Type=oneshot
ExecStart=/usr/bin/mcsrv-script timer_two
