# srv/mcsrv/.config/systemd/user/mcsrv-timer-1.service
[Unit]
Description=McSrv Script Timer 1 Service

[Service]
Type=oneshot
ExecStart=/usr/bin/mcsrv-script timer_one
